# sys/usr/lib/systemd/user/kitty.service
[Unit]
Description=GPU based terminal emulator (headless single instance)
Documentation=https://sw.kovidgoyal.net/kitty/
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
ExecStart=/usr/bin/kitty --single-instance --start-as=hidden
Restart=on-failure

[Install]
WantedBy=graphical-session.target
